# supervisord.conf
[supervisord]
nodaemon=true

[program:sickrage]
user=sickrage
command=/opt/sickrage/SickBeard.py --nolaunch --pidfile=/home/sickrage/.sickrage.pid --datadir=/config
